# deploy/centos7/material-system.service
[Unit]
Description=Material Management System
After=network.target

[Service]
Type=simple
User=django
Group=django
WorkingDirectory=/home/django/material_system
Environment=PATH=/home/django/.local/bin:/usr/local/bin:/usr/bin
Environment=DJANGO_SETTINGS_MODULE=material_system.settings
Environment=DEBUG=False
ExecStart=/home/django/.local/bin/gunicorn material_system.wsgi:application --bind 0.0.0.0:8000 --workers 3 --timeout 120 --keep-alive 5
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=10

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/django/material_system

[Install]
WantedBy=multi-user.target